# docker/local/django/crontab
*/3 * * * * /usr/src/app/manage.py example
